# systemd/user/gpu-fan.timer
[Unit]
Description=Timer for NVIDIA fan speed adjuster

[Timer]
AccuracySec=1s
OnActiveSec=5s
OnUnitActiveSec=5sec

[Install]
WantedBy=timers.target
