# systemd/user/gpu-fan.timer
[Unit]
Description=Timer for NVIDIA fan speed adjuster

[Timer]
OnBootSec=1min
OnUnitActiveSec=10sec

[Install]
WantedBy=timers.target
